# deploy/systemd/secrets-server.service
[Unit]
Description=Secrets Manager server
Documentation=https://github.com/darui3018823/Secrets-Manager
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
User=secrets
Group=secrets
WorkingDirectory=/var/lib/secrets-manager

Environment=SECRETS_DB_PATH=/var/lib/secrets-manager/secrets.db
Environment=SECRETS_AUDIT_PATH=/var/log/secrets-manager/audit.jsonl
Environment=SECRETS_BIND=127.0.0.1:8787
LoadCredential=secrets-passphrase:/etc/secrets-manager/master-passphrase

ExecStart=/usr/local/bin/secrets-server serve
Restart=on-failure
RestartSec=5s
UMask=0077

StateDirectory=secrets-manager
LogsDirectory=secrets-manager
RuntimeDirectory=secrets-manager

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/secrets-manager /var/log/secrets-manager
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
